# supervisord.conf
[supervisord]
nodaemon=true

[program:startup_run]
command=sh ~/run.sh




